# config/.config/systemd/user/xdg-desktop-portal-wlr.service.d/override.conf
[Unit]
After=sway-session.target
PartOf=sway-session.target
