# config/.config/systemd/user/waybar.service.d/override.conf
[Unit]
PartOf=sway-session.target
After=sway-session.target
